# crontab
0,15,30,45 * * * * /workspace/bin/rails airbrake:test

# An empty line is required at the end of this file for a valid cron file.
